# systemd/minusculeclaw-webhook.service
[Unit]
Description=MinusculeClaw Webhook Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/minusculeclaw
EnvironmentFile=%h/minusculeclaw/.env
ExecStart=/usr/bin/python3 %h/minusculeclaw/webhook_server.py
Restart=always
RestartSec=3

[Install]
WantedBy=default.target
