# systemd/shellclaw-heartbeat.service
[Unit]
Description=ShellClaw Daily Heartbeat

[Service]
Type=oneshot
WorkingDirectory=%h/shellclaw
EnvironmentFile=%h/shellclaw/.env
ExecStart=/bin/bash -lc 'cd %h/shellclaw && ./scripts/heartbeat.sh'
